# deploy/systemd/slt-setup.service
# slt-setup.service
#
# Privileged one-shot that prepares networking for slt-server: creates and
# configures the TUN device, enables IPv4 forwarding, and installs NAT. Runs as
# root. Its ExecStop reverses everything, so cleanup happens whenever the server
# is stopped, restarted, or the system shuts down (see PartOf= below).

[Unit]
Description=SLT VPN networking setup (TUN, IPv4 forwarding, NAT)
Wants=network-online.target
After=network-online.target
Before=slt-server.service
# Follow slt-server's stop/restart so teardown runs automatically when the
# server is stopped or restarted (and on system shutdown).
PartOf=slt-server.service

[Service]
Type=oneshot
RemainAfterExit=yes
# TUN interface name/address/prefix/MTU are read from SLT_CONFIG's [tun]
# section by `slt net`. Only deployment-local values live here.
Environment=SLT_CONFIG=/etc/slt/server.toml SLT_USER=slt SLT_GROUP=slt
ExecStart=/usr/local/bin/slt net up --config ${SLT_CONFIG} --user ${SLT_USER} --group ${SLT_GROUP} --ipv4-forward --masquerade
ExecStop=/usr/local/bin/slt net down --config ${SLT_CONFIG} --masquerade

[Install]
WantedBy=multi-user.target
